# deploy/steam-group-rec.service
[Unit]
Description=Steam Group Recommender
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/steam-group-rec
# 只监听本机，公网流量必须走 nginx 的 HTTPS 反代，避免 key 明文暴露在 8673。
Environment=STEAMREC_HOST=127.0.0.1
Environment=STEAMREC_PORT=8673
Environment=STEAMREC_DATA_DIR=/opt/steam-group-rec/data
# 只承载访问统计令牌；AI key 由用户网页填写，不落服务器。
EnvironmentFile=-/etc/steam-group-rec.env
ExecStart=/opt/steam-group-rec/.venv/bin/python /opt/steam-group-rec/app.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
